# srcs/influxDB/supervisor.conf
[supervisord]
nodaemon=true

[program:influxd]
command=sh -c "influxd && pkill supervisord"

[program:telegraf]
command=sh -c "/etc/telegraf/usr/bin/telegraf && pkill supervisord"