# srcs/phpMyAdmin/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=sh -c "nginx -g 'daemon off;' && pkill supervisord"

[program:php-fpm7]
command=sh -c "php-fpm7 -F -R && pkill supervisord"